# rootfs-custom/mdpart2dm/lib/udev/rules.d/01-mdpart-dm-unmap.rules
KERNEL=="md*p*", ACTION=="remove", SUBSYSTEM=="block", RUN+="mdpart2dm --remove $kernel"
