# rootfs-custom/mdpart2dm/lib/udev/rules.d/99-mdpart-dm-map.rules
KERNEL=="md*p*", ACTION=="add", SUBSYSTEM=="block", RUN+="mdpart2dm --create $kernel"
