# rootfs/opt/rc.local
#!/bin/sh

export PATH=/usr/local/bin:$PATH

# *** JW-NAS-CONF-START ***
# *** JW-NAS-CONF-END ***
