# rootfs/lib/udev/rules.d/99-sd-remove-stick.rules
ACTION=="add", KERNEL=="sd*|md*", RUN+="/bin/chmod -t %N >/tmp/a 2>/tmp/b"
