# rootfs-debian7/etc/cron.d/logrotate
# execute logrotate per 5 minutes
*/5 * * * *	root	[ -x /usr/sbin/logrotate ] && sudo /usr/sbin/logrotate /etc/logrotate.conf
